# udev/01-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/home/joe/.resources/udev-resize.sh"
